# gzmo-daemon/gzmo-daemon.service.template
[Unit]
Description=GZMO Daemon — Sovereign AI Heartbeat
After=network.target

[Service]
Type=simple
User=%u
# This file is a TEMPLATE. Install it using `./install_service.sh`.
#
# The installer replaces:
# - __GZMO_DAEMON_DIR__ (absolute path to `gzmo-daemon/`)
# - __BUN_BIN__ (bun executable path)
WorkingDirectory=__GZMO_DAEMON_DIR__
EnvironmentFile=__GZMO_DAEMON_DIR__/.env
ExecStart=__BUN_BIN__ run index.ts
Restart=on-failure
RestartSec=10

# Ollama performance tuning (if Ollama runs on same host)
Environment=OLLAMA_KV_CACHE_TYPE=q8_0
Environment=OLLAMA_FLASH_ATTENTION=1
Environment=OLLAMA_KEEP_ALIVE=-1

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=gzmo-daemon

[Install]
WantedBy=default.target
