# gzmo-daemon/gzmo-daemon.service
[Unit]
Description=GZMO Daemon — Sovereign AI Heartbeat
After=network.target

[Service]
Type=simple
# ── Paths: adjust to match your system layout ──
# The defaults below assume edge-node is at ~/Dokumente/Playground/DevStack_v2/edge-node
User=%u
WorkingDirectory=%h/Dokumente/Playground/DevStack_v2/edge-node/gzmo-daemon
EnvironmentFile=%h/Dokumente/Playground/DevStack_v2/edge-node/gzmo-daemon/.env
ExecStart=%h/.bun/bin/bun run index.ts
Restart=on-failure
RestartSec=10

# Ollama performance tuning (if Ollama runs on same host)
Environment=OLLAMA_KV_CACHE_TYPE=q8_0
Environment=OLLAMA_FLASH_ATTENTION=1
Environment=OLLAMA_KEEP_ALIVE=-1

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=gzmo-daemon

[Install]
WantedBy=default.target
